# .config/systemd/user/change_wallpaper.timer
[Unit]
Description=Timer to change wallpapers every hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
