# .config/systemd/user/change_wallpaper.service
[Unit]
Description=Change wallpaper using SWWW

[Service]
ExecStart=/bin/bash %h/.config/awww/change_wallpaper.sh

[Install]
WantedBy=default.target
