# .config/systemd/user/change_wallpaper.timer
[Unit]
Description=Timer to change wallpapers every hour

[Timer]
# Every 5mins
OnCalendar=*:0/5
Persistent=true

[Install]
WantedBy=timers.target
